# conf.d/wmaker.conf
[program:wmaker]
command=wmaker
autorestart=true
